# setup_ansible/pine/files/mqttwarn.supervisord.conf
[program:mqttwarn]
command = /usr/local/bin/mqttwarn
user = ubuntu
environment= MQTTWARNINI="/home/ubuntu/mqttwarn.ini"